# supervisord.conf
[supervisord]
nodaemon=true

[program:puma]
directory=/data
command=/usr/local/bin/bundle exec puma -C /data/config/puma.rb -e production
autorestart=true
